# profile.d/python.profile
export WORKON_HOME=$HOME/.virtualenvs
export VIRTUALENVWRAPPER_PYTHON=/usr/bin/python3
export PROJECT_HOME=$HOME/ws
source $HOME/.local/bin/virtualenvwrapper.sh

# zsh pyenv plugin does that already
# export PATH="/home/behaghel/.pyenv/bin:$PATH"
# eval "$(pyenv init -)"
# eval "$(pyenv virtualenv-init -)"
